# Startup/rc.local
#Add the following lines
sleep 20

sudo /home/pi/AutoCam.sh 